# supervisord.conf
[supervisord]
nodaemon=true
logfile=/dev/stdout
logfile_maxbytes=0
loglevel=info
pidfile=/var/run/supervisord.pid

[unix_http_server]
file=/var/run/supervisor.sock

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

# Nginx - 反向代理（延迟启动，等待Next.js就绪）
[program:nginx]
command=nginx -g "daemon off;"
autostart=true
autorestart=true
startretries=3
startsecs=10
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=300

# Next.js应用 - Web服务器（优先启动）
[program:nextjs]
command=node server.js
directory=/app
autostart=true
autorestart=true
startretries=3
startsecs=5
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=100
user=nextjs
environment=NODE_ENV="production",PORT="3000",HOSTNAME="0.0.0.0",NODE_OPTIONS="--max-old-space-size=4096"

# 调度器 - 定时任务
[program:scheduler]
command=node dist/scheduler.js
directory=/app
autostart=true
autorestart=true
startretries=3
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=300
user=nextjs
environment=NODE_ENV="production",TZ="Asia/Shanghai",NODE_OPTIONS="--max-old-space-size=4096"

# Python Google Ads Service - 服务账号模式API
[program:python-ads-service]
command=python3 -m uvicorn main:app --host 0.0.0.0 --port 8001
directory=/app/python-service
autostart=true
autorestart=true
startretries=3
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=200
user=nextjs
